# container-qdrant.service
# container-qdrant.service
# autogenerated by Podman 5.4.0
# Wed Aug 27 11:53:45 EDT 2025

[Unit]
Description=Podman container-qdrant.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman start qdrant
ExecStop=/usr/bin/podman stop  \
	-t 10 qdrant
ExecStopPost=/usr/bin/podman stop  \
	-t 10 qdrant
PIDFile=/run/containers/storage/overlay-containers/1f35df64124baf8b7a41a18e3fd72ca9610b8d1bd44ead0b6173b51660868eab/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
